# 1.2.4/distributed/supervisord.conf
[supervisord]
nodaemon=true

[program:tomcat7]
command=/usr/local/tomcat/bin/catalina.sh run
environment=CATALINA_BASE=/usr/local/tomcat,CATALINA_TMPDIR=/usr/local/tomcat/tmp
